# contrib/openwrt/ubd/files/ubd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=80

BIN=ubd
DEFAULT=/etc/default/$BIN
OPTIONS=""

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	$BIN $OPTIONS
}

stop() {
  killall $BIN
}


