# contrib/openwrt/ubd/files/ubnetd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=80

BIN=ubnetd
DEFAULT=/etc/default/$BIN
INTERFACE="eth0"

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	$BIN $INTERFACE
}

stop() {
  killall $BIN
}


